# systemd/docker-prune.service
[Unit]
Description=Prune unused Docker resources

[Service]
ExecStart=/usr/bin/docker system prune -f